# java/src/deb/init/monasca-persister.conf
# Startup script for the monasca-persister

description "Monasca Persister java app"
start on runlevel [2345]

console log
respawn

setgid monasca
setuid persister
exec /usr/bin/java -Dfile.encoding=UTF-8 -Xmx8g -cp /opt/monasca/monasca-persister.jar:/opt/monasca/vertica/vertica_jdbc.jar monasca.persister.MonPersisterApplication server /etc/monasca/persister-config.yml
